# bolixoserv.sysv
#!/bin/sh
#chkconfig: 30 99 09

case "$1" in
start)
	if [ -x /root/bolixostart.sh ]; then
		/root/bolixostart.sh
	else
		echo /root/bolixostart.sh missing
		echo run bolixo-production config
	fi
	;;
stop)
	/root/bolixostop.sh
	;;
*)
	echo start or stup
	;;
esac

